# MODS
1. /boot/commandline added for splashscreen:
quiet

modified for PiTFT in portrait mode:
fbtft_device.rotate=0

2. added file for splashscreen made executable using "sudo chmod a+x /etc/init.d/asplashscreen":
/etc/init.d/asplashscreen

#! /bin/sh
### BEGIN INIT INFO
# Provides:          asplashscreen
# Required-Start:
# Required-Stop:
# Should-Start:      
# Default-Start:     S
# Default-Stop:
# Short-Description: Show custom splashscreen
# Description:       Show custom splashscreen
### END INIT INFO

do_start () {

    omxplayer /ect/occ/occ_splashscreen.mp4 &  
    exit 0
}

case "$1" in
  start|"")
    do_start
    ;;
  restart|reload|force-reload)
    echo "Error: argument '$1' not supported" >&2
    exit 3
    ;;
  stop)
    # No-op
    ;;
  status)
    exit 0
    ;;
  *)
    echo "Usage: asplashscreen [start|stop]" >&2
    exit 3
    ;;
esac

:

3. added service:
sudo insserv /etc/init.d/asplashscreen

4. Added splashscreen video:
/etc/occ/occ_splashscreen.mp4: ISO Media, MPEG v4 system, version 1

  Duration: 00:00:20.13, start: 0.000000, bitrate: 209 kb/s
    Stream #0:0(und): Video: h264 (High) (avc1 / 0x31637661), yuv420p, 240x320 [SAR 1:1 DAR 3:4], 11 kb/s, 19.98 fps, 19.98 tbr, 20k tbn, 39.96 tbc (default)

5. Install mplayer
6. Add fbcp and make /etc/init.d/asplashscreen use it. fbcp: https://github.com/tasanakorn/rpi-fbcp

#! /bin/sh
### BEGIN INIT INFO
# Provides:          asplashscreen
# Required-Start:
# Required-Stop:
# Should-Start:      
# Default-Start:     S
# Default-Stop:
# Short-Description: Show custom splashscreen
# Description:       Show custom splashscreen
### END INIT INFO

do_start () {
	/usr/bin/occ_splashscreen &
	exit 0
}

case "$1" in
  start|"")
    do_start
    ;;
  restart|reload|force-reload)
    echo "Error: argument '$1' not supported" >&2
    exit 3
    ;;
  stop)
    # No-op
    ;;
  status)
    exit 0
    ;;
  *)
    echo "Usage: asplashscreen [start|stop]" >&2
    exit 3
    ;;
esac

:


7. Add service to runlevel 2 
cd /etc/rc2.d
sudo ln -s ../init.d/asplashscreen S01asplashscreen

8. sudo apt-get remove wolfram-engine
9. Install font Segment14 in ~/.fonts
http://openfontlibrary.org/en/font/segment14#Segment14-Regular
10. Calibrated touchscreen using 
sudo TSLIB_FBDEVICE=/dev/fb1 TSLIB_TSDEVICE=/dev/input/touchscreen ts_calibrate

and test using:
sudo TSLIB_FBDEVICE=/dev/fb1 TSLIB_TSDEVICE=/dev/input/touchscreen ts_test

10. sudo apt-get install python3-lxml


11. sudo apt-get install minicom

12. Remove tty from /boot/cmdline.txt
13. Comment out serial console from /etc/inittab
14. sudo apt-get install gpsd gpsd-clients python-gps

15. Reconfigure gpsd config file to run gpsd by default
16. sudo apt-get install python3-tz
17. Install tzwhere for python from https://github.com/pegler/pytzwhere
18. sudo apt-get install python-quantities
19. Add /usr/bin/occ:

#!/bin/bash
cd /home/pi/OpenCyclingComputer/
sudo python src/occ.py
20. Add service to runlevel 2
cd /etc/rc2.d
sudo ln -s ../init.d/occ S99occ
21. Add /etc/init.d/occ:

#! /bin/sh
### BEGIN INIT INFO
# Provides:          OpenCyclingComputer
# Required-Start:
# Required-Stop:
# Should-Start:      
# Default-Start:     S
# Default-Stop:
# Short-Description: Open Cycling Computer
# Description:       Open Cycling Computer
### END INIT INFO

do_start () {
	/usr/bin/occ &
	exit 0
}

case "$1" in
  start|"")
    do_start
    ;;
  restart|reload|force-reload)
    echo "Error: argument '$1' not supported" >&2
    exit 3
    ;;
  stop)
    # No-op
    ;;
  status)
    exit 0
    ;;
  *)
    echo "Usage: occ [start|stop]" >&2
    exit 3
    ;;
esac

:
